# overlay/etc/systemd/system/server.service
[Unit]
Description=server
Requires=ttyACM0.device
After=systemd-modules-load.service

[Service]
User=root
Group=root
Type=simple 
WorkingDirectory=/scripts
ExecStart=python3 server.py &
Restart=always
RestartSec=10 

[Install]
WantedBy=multi-user.target
